# launchd/com.butler.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.butler.agent</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Absolute path to venv python — updated by install.sh -->
        <string>BUTLER_HOME/.venv/bin/python</string>
        <string>-m</string>
        <string>butler.main</string>
        <string>BUTLER_HOME/config/butler.yaml</string>
    </array>

    <key>WorkingDirectory</key>
    <string>BUTLER_HOME</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin:/usr/local/opt/node/bin</string>
        <key>PLAYWRIGHT_BROWSERS_PATH</key>
        <string>BUTLER_HOME/.playwright</string>
        <key>HOME</key>
        <string>HOME_DIR</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>BUTLER_HOME/logs/butler.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>BUTLER_HOME/logs/butler.stderr.log</string>
</dict>
</plist>
